# init/prefdm.conf
# prefdm - preferred display manager
#
# Starts gdm/xdm/etc by preference
#
# Do not edit this file directly. If you want to change the behaviour,
# please create a file prefdm.override and put your changes there.

start on stopped rc RUNLEVEL=5

stop on starting rc RUNLEVEL=[!5]

console output
respawn
respawn limit 10 120
exec /etc/X11/prefdm -nodaemon
