# init/reload-init.conf
start on started messagebus

task
exec /bin/kill -HUP 1
