# udev/rules.d/55-ccw.rules
ACTION=="add", SUBSYSTEM=="ccw", RUN+="/lib/udev/ccw_init"
