# udev/rules.d/60-net.rules
SUBSYSTEM=="net", RUN+="/etc/sysconfig/network-scripts/net.hotplug"
