# systemd/system/brandbot.service
[Unit]
Description=Flexible Branding Service

[Service]
ExecStart=/usr/sbin/brandbot
PrivateTmp=true
