# systemd/system/brandbot.service
[Unit]
Description=Flexible Branding Service

[Service]
ExecStart=/usr/libexec/initscripts/brandbot
PrivateTmp=true
