# etc/NetworkManager/dispatcher.d/00-netreport
#!/bin/sh

cd /etc/sysconfig/network-scripts
. ./network-functions

[ -f ../network ] && . ../network

# Notify programs that have requested notification
do_netreport

exit 0
